# cronjobs/crontab
# coinmines crontab


* * * * *	/data/coinmines/cronjobs/cronjob.sh
